# package/debian/postinst
#!/bin/bash

case "$1" in
	configure)
		systemctl enable --now snx-rs 2>/dev/null >/dev/null || true
		;;
	*)
		;;
esac

exit 0
