# assets/snx-rs.service
[Unit]
Description=SNX-RS VPN client
After = network.target network-online.target
Wants = network-online.target

[Service]
Type=simple
ExecStart=/opt/snx-rs/snx-rs -m command -l info

[Install]
WantedBy=multi-user.target
